# scripts/bias_scheduler.service
[Unit]
Description=HKJC Bias Optimization Scheduler
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/hkjc
Environment=PATH=/opt/hkjc/.venv/bin
ExecStart=/opt/hkjc/.venv/bin/python /opt/hkjc/scripts/start_bias_scheduler.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
